# debian/etc/init/monaco-stats.conf
description "Stats reporting for Monaco"
author "keith@hulu.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid monaco
setgid monaco
exec /usr/bin/python /usr/monaco/monaco_stats.py
